# supporting_files/supervisord-memcache.conf
[program:memcached]
command=/start-memcache.sh
numprocs=1
autostart=true
autorestart=true
